# system/sobu7.cron
#!/bin/sh

ps -elf |grep -v grep | grep 'calling-sobu46'

if [[ $? -ne 0 ]]
then
	/root/netcat.sobu7-calling-sobu46-on-port-4455.sh &
fi

ps -elf |grep -v grep | grep 'listening-to-10.0.1.201'
if [[ $? -ne 0 ]]
then
	/root/netcat.sobu7-listening-to-10.0.1.201-on-port-4444.sh &
fi
