# system/tscp.cron
#!/bin/bash

ps -elf |grep -v grep | grep 'receive_urms_data'

if [[ $? -ne 0 ]]
then
	/home/art_coord_ramp_metering/system/artgo
fi
