# system/urms.cron
* * * * * /home/art_coord_ramp_metering/system/check_and_restart_urms.sh
